# branches/per-image-credit/system/webserver/etc/cron.daily/synctime
#!/bin/bash
/usr/sbin/ntpdate ntp.demon.co.uk

